# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:bot]
command=python bot.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/bot.err.log
stdout_logfile=/var/log/bot.out.log
environment=TARGET_CHANNEL="%(ENV_TARGET_CHANNEL)s",CHECK_COOLDOWN="%(ENV_CHECK_COOLDOWN)s",PANEL_USERNAME="%(ENV_PANEL_USERNAME)s",PANEL_PASS="%(ENV_PANEL_PASS)s",PANEL_ADDRESS="%(ENV_PANEL_ADDRESS)s",BOT_TOKEN="%(ENV_BOT_TOKEN)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",PROXY_DOMAIN="%(ENV_PROXY_DOMAIN)s",PROXY_PORT="%(ENV_PROXY_PORT)s",ADMIN_ID="%(ENV_ADMIN_ID)s"

[program:proxy]
command=python proxy_server.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/proxy.err.log
stdout_logfile=/var/log/proxy.out.log
environment=PROXY_DOMAIN="%(ENV_PROXY_DOMAIN)s",PROXY_PORT="%(ENV_PROXY_PORT)s" 