# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:bot]
command=python bot.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/bot.err.log
stdout_logfile=/var/log/bot.out.log

[program:proxy]
command=python proxy_server.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/proxy.err.log
stdout_logfile=/var/log/proxy.out.log 